# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

# Control socket so a fatal failure (e.g. the S3 mount) can shut the whole
# container down instead of leaving an FTP server running over an empty bucket.
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# One-shot: mounts the bucket and provisions users, then exits 0. If it fails it
# shuts supervisord down (see s3-fuse.sh), so the container exits and can be
# restarted cleanly rather than serving an empty mount point.
[program:s3-fuse]
command=/usr/local/s3-fuse.sh
autorestart=false
startsecs=0
priority=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

# Waits for the S3 mount to be ready before serving (start-vsftpd.sh).
[program:vsftpd]
command=/usr/local/start-vsftpd.sh
autorestart=true
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:add_users_in_container]
command=/usr/local/add_users_in_container.sh
autorestart=true
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
